# rootfs.centos.8/etc/supervisord.d/novnc.ini
[supervisord]
autorestart=true
directory=/home/user
nodaemon=true
redirect_stderr=true
stopsignal=QUIT
user=user

[group:application]
programs=vnc,novnc,xfce

[program:vnc]
command=/usr/bin/Xvnc -localhost no -SecurityTypes None
priority=10
stopsignal=KILL

[program:novnc]
command=/opt/novnc/utils/novnc_proxy
directory=/opt/novnc
priority=20
stopasgroup=true

[program:xfce]
command=/usr/bin/startxfce4
directory=/home/user
priority=30
